# deploy/supervisor/gekychat-worker.conf
; GekyChat Laravel queue worker – Supervisor config
; Server: gekymedia.com
; App path: /home/gekymedia/web/chat.gekychat.com/public_html
;
; One-time setup on server:
;   1. Install supervisor (if not installed):
;      apt-get install supervisor   # Debian/Ubuntu
;   2. Copy this file to /etc/supervisor/conf.d/gekychat-worker.conf
;      (or symlink from repo: ln -s /home/gekymedia/web/chat.gekychat.com/public_html/deploy/supervisor/gekychat-worker.conf /etc/supervisor/conf.d/)
;   3. Create log directory: mkdir -p /home/gekymedia/web/chat.gekychat.com/public_html/storage/logs
;   4. sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start gekychat-worker:*
;
; After deploy, deploy.ps1 runs queue:restart so workers pick up new code.

[program:gekychat-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/gekymedia/web/chat.gekychat.com/public_html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=gekymedia
numprocs=2
redirect_stderr=true
stdout_logfile=/home/gekymedia/web/chat.gekychat.com/public_html/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
