# deploy/supervisor/gekychat-reverb.conf
; GekyChat Laravel Reverb WebSocket Server – Supervisor config
; Server: gekymedia.com
; App path: /home/gekymedia/web/chat.gekychat.com/public_html
;
; One-time setup on server:
;   1. Copy/symlink this file to /etc/supervisor/conf.d/gekychat-reverb.conf
;   2. sudo supervisorctl reread && sudo supervisorctl update
;   3. sudo supervisorctl start gekychat-reverb:*
;
; Reverb runs on port 8080 by default, nginx proxies /app and /apps to it.

[program:gekychat-reverb]
process_name=%(program_name)s
command=php /home/gekymedia/web/chat.gekychat.com/public_html/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=gekymedia
numprocs=1
redirect_stderr=true
stdout_logfile=/home/gekymedia/web/chat.gekychat.com/public_html/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
